# divers/crontab
*/5 * * * * /bin/bash $(cd /srv/web/devfly.fr/ && git pull)
